# .config/systemd/user/borgmatic.timer
[Unit]
Description=Weekly Backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
